# automation/com.khutbahnotes.weekly-masjid-publishing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.khutbahnotes.weekly-masjid-publishing</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/abbasanwar/Desktop/Projects/Khutbah Notes AI/scripts/run_weekly_masjid_workflow.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/Users/abbasanwar/Desktop/Projects/Khutbah Notes AI</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Weekday</key>
    <integer>7</integer>
    <key>Hour</key>
    <integer>20</integer>
    <key>Minute</key>
    <integer>40</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/Users/abbasanwar/Desktop/Projects/Khutbah Notes AI/logs/weekly-masjid-publishing.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/abbasanwar/Desktop/Projects/Khutbah Notes AI/logs/weekly-masjid-publishing.log</string>
</dict>
</plist>
